# crontab
* * * * * root /srv/mixpanel-ddc/clear_logs.sh
# * * * * * root /srv/mixpanel-ddc/restart_full.sh
