# crontab
* * * * * root /srv/mixpanel-ddc/clear_logs.sh
